# debian/cron.d
#
# Regular cron jobs for the mirthless package
#
0 4 * * *  root  mirthless_maintenance
